# board/opencs18/rootfs_overlay/etc/network/if-pre-up.d/eth0-hwaddr
#!/bin/sh

HWAFILE=/etc/network/hwaddr.txt

if [ "$IFACE" = "eth0" ]; then
    if [ -f $HWAFILE ]; then
	HWADDR=`cat $HWAFILE`
	ifconfig $IFACE hw ether $HWADDR
    fi
fi